# installer/autostart-terminal.service
[Unit]
Description=Auto start terminal on boot

[Service]
Type=simple
ExecStart=/usr/bin/gnome-terminal

[Install]
WantedBy=multi-user.target